# deploy/certbot-renew.service
# ===============================================
# 文件路径：deploy/certbot-renew.service
# 功能：定义 Certbot 自动续签服务
# ===============================================

[Unit]
Description=Renew Let's Encrypt certificates using Certbot
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/certbot-renew.sh
# 环境变量（确保 systemd 能找到 certbot 命令）
Environment="PATH=/usr/bin:/usr/local/bin"
