# deploy/certbot-renew.timer
# ===============================================
# 文件路径：deploy/certbot-renew.timer
# 功能：每天凌晨定时执行 Certbot 自动续签
# ===============================================

[Unit]
Description=Run Certbot renew daily

[Timer]
# 每天凌晨 3:00:00 触发
OnCalendar=*-*-* 03:00:00
# Persistent=true 表示即使服务器重启，错过的任务也会补执行
Persistent=true

[Install]
WantedBy=timers.target
